# outputs/adshares-node-docker/supervisord.conf
; Compatibility reference only. Docker Compose starts adsd directly and does
; not install or use Supervisor.
[program:adsd]
command=/usr/local/bin/adsd --init=true -w /home/ads/.adsd
directory=/home/ads/.adsd
user=ads
autostart=true
autorestart=true
startsecs=10
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
